# trans_rights_bot/docker_runtime/bot.crontab
0 10 * * * /usr/local/bin/trans_rights_bot